# buildroot-external/board/common/rootfs-overlay/etc/systemd/system/sshd.service.d/override.conf
[Service]
ExecStartPre=
ExecStartPre=/usr/bin/ssh-keygen -A -f /mnt/overlay
