# buildroot/board/neuron-rpi3-64/rootfs_overlay/usr/lib/systemd/system/mnt-data.mount
[Unit]
Description=Data partition
DefaultDependencies=no
Before=umount.target local-fs.target
Conflicts=umount.target
After=systemd-fsck@dev-disk-by\x2dlabel-\x2ddata.service
Wants=systemd-fsck@dev-disk-by\x2dlabel-\x2ddata.service

[Mount]
What=/dev/disk/by-label/data
Where=/mnt/data
Type=ext4
Options=defaults,noatime,commit=30

[Install]
WantedBy=local-fs.target
