# buildroot-external/board/common/rootfs-overlay/etc/systemd/system/owhttpd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/bin/owhttpd --foreground --server=127.0.0.1 -p 8080
User=root
Group=root
