# buildroot/package/unipi-tools/usr/lib/systemd/system/unipiconfig.service
[Unit]
Description=UniPi Config

[Service]
ExecStart=/opt/unipi/tools/unipiconfig.sh

[Install]
RequiredBy=unipitcp.service
