# buildroot-external/board/common/rootfs-overlay/etc/systemd/system/systemd-firstboot.service.d/override.conf
[Unit]
After=etc.mount

[Service]
ExecStart=systemd-firstboot --root=/ --setup-machine-id
