# buildroot-external/board/common/rootfs-overlay/etc/profile.d/python.sh
#!/bin/sh

export PYTHONUSERBASE=/usr/local/
