# buildroot-external/board/common/rootfs-overlay/etc/profile.d/python.sh
export PYTHONUSERBASE=/usr/local
export PATH="${PATH}:/usr/local/opt/venv/unipi-control/bin"
