# buildroot-external/board/common/rootfs-overlay/etc/profile.d/python.sh
#!/bin/dash

# Set python user base to /usr/local.
# pip install packages without venv to this location.
export PYTHONUSERBASE=/usr/local
