# buildroot-external/board/common/rootfs-overlay/etc/profile.d/unipi-control.sh
#!/bin/sh

export PATH="${PATH}:/usr/local/opt/venv/unipi-control/bin"

/usr/local/bin/os-tools check-version unipi-control
