# buildroot-external/board/common/rootfs-overlay/etc/profile.d/unipi-control.sh
#!/bin/sh

export PATH="/mnt/data/venv/unipi-control/bin:${PATH}"

/usr/bin/check-unipi-control-version.sh
