# buildroot-external/board/common/rootfs-overlay/usr/lib/udev/rules.d/90-unipi-id.rules
SUBSYSTEM=="platform", ACTION=="bind", DEVPATH=="/devices/platform/unipi-id", TAG+="systemd", RUN+="/opt/unipi/tools/uhelper unipi-id unipi-id"
