# buildroot-external/board/common/rootfs-overlay/etc/systemd/system/docker.service.d/override.conf
[Unit]
RequiresMountsFor=/etc/docker /mnt/data

[Service]
OOMScoreAdjust=-400
