# buildroot-external/board/common/rootfs-overlay/etc/systemd/system/owserver.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/owserver --i2c=/dev/i2c-1:ALL
User=root
Group=root
