# buildroot-external/board/common/rootfs-overlay/usr/lib/systemd/system/machine-id-setup.service
[Unit]
Description=Set up machine-id
DefaultDependencies=no
After=etc.mount
ConditionPathIsReadWrite=/etc
ConditionFileNotEmpty=!/etc/machine-id

[Service]
Type=simple
ExecStart=systemd-machine-id-setup
