# buildroot-external/board/common/rootfs-overlay/usr/lib/systemd/system/update-hostname.service
[Unit]
Description=Update hostname
Before=systemd-networkd.target
After=etc.mount
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/hostname -F /etc/hostname

[Install]
WantedBy=multi-user.target
