# buildroot-external/board/common/rootfs-overlay/usr/lib/systemd/system/var-log.mount
[Unit]
Description=Persistent /var/log
Requires=mnt-data.mount
After=mnt-data.mount
Before=docker.service mosquitto.service

[Mount]
What=/mnt/data/logs
Where=/var/log
Type=none
Options=bind

[Install]
WantedBy=os-bind.target
