# buildroot-external/board/common/rootfs-overlay/etc/systemd/system/systemd-networkd.service.d/override.conf
[Unit]
RequiresMountsFor=/etc/hostname /etc/hosts
